# scripts/teamkim-bml1.service
[Unit]
Description=teamkim-bml1 maze robot — run main.py on boot
After=multi-user.target

[Service]
Type=simple
User=team2
Group=team2
WorkingDirectory=/home/team2/teamkim-bml1
ExecStart=/bin/bash /home/team2/teamkim-bml1/scripts/run_on_boot.sh

# One drive per power-cycle. Don't auto-restart — if main.py exits we want
# the car to STAY stopped, not start a fresh run on its own.
Restart=no
TimeoutStopSec=2

# Send stdout/stderr to the journal so `journalctl -u teamkim-bml1 -b`
# shows what happened this boot.
StandardOutput=journal
StandardError=journal

# Tunables. Override per-instance with: sudo systemctl edit teamkim-bml1
# (creates a drop-in at /etc/systemd/system/teamkim-bml1.service.d/override.conf)
Environment=REPO_DIR=/home/team2/teamkim-bml1
Environment=BOOT_DELAY_S=10
Environment=RUN_DURATION_S=300
Environment=RUN_NAME=boot_run
# Full wall-following maze loop. Override to integration_test.py only for
# a deliberate non-maze hardware scenario.
Environment=RUN_SCRIPT=main.py

[Install]
WantedBy=multi-user.target
